# packages/google-compute-engine/src/etc/init/google-accounts-daemon.conf
# Manages accounts from metadata SSH keys.
start on started google-network-daemon

respawn
exec /usr/bin/google_accounts_daemon
