# packages/google-compute-engine/src/etc/init/google-instance-setup.conf
# Runs instance setup on boot.
start on started rsyslog

task

exec /usr/bin/google_instance_setup
