# packages/google-compute-engine/src/etc/init/google-network-daemon.conf
# Manages network interfaces.
start on stopped google-instance-setup

respawn
exec /usr/bin/google_network_daemon
