# packages/google-compute-engine/src/lib/systemd/system/google-clock-skew-daemon.service
[Unit]
Description=Google Compute Engine Clock Skew Daemon
After=network.target google-instance-setup.service google-network-daemon.service

[Service]
Type=simple
ExecStart=/usr/bin/google_clock_skew_daemon
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
