# packages/google-compute-engine/src/lib/systemd/system/google-network-daemon.service
[Unit]
Description=Google Compute Engine Network Daemon
After=network-online.target network.target
After=google-instance-setup.service
PartOf=network.service

[Service]
Type=simple
ExecStart=/usr/bin/google_network_daemon
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
