# packages/google-compute-engine/src/lib/systemd/system/google-instance-setup.service
[Unit]
Description=Google Compute Engine Instance Setup
After=network-online.target network.target rsyslog.service
Before=sshd.service

[Service]
Type=oneshot
ExecStart=/usr/bin/google_instance_setup
StandardOutput=journal+console

[Install]
WantedBy=sshd.service
WantedBy=multi-user.target
