# google_compute_engine_init/upstart/google-startup-scripts.conf
# Runs a startup script from metadata.
start on started google-network-deamon

exec /usr/bin/google_metadata_script_runner --script-type startup
