# google_compute_engine_init/upstart/google-accounts-daemon.conf
# Manages accounts from metadata SSH keys.
start on stopped google-network-setup

respawn
exec /usr/bin/google_accounts_daemon
