# google_compute_engine_init/upstart/google-network-setup.conf
# Enables network interfaces on boot.
start on stopped google-instance-setup

task

exec /usr/bin/google_network_setup
