# google_compute_engine_init/upstart/google-clock-skew-daemon.conf
# Sync the system clock on migration.
start on started google-network-daemon

respawn
exec /usr/bin/google_clock_skew_daemon
