# google-daemon/etc/init/google-address-manager.conf
# This service configures local addresses in a Google Compute Engine instance.
start on google-rc-local-has-run

respawn
exec /usr/share/google/google_daemon/manage_addresses.py
