# google-daemon/etc/init/google-accounts-manager-task.conf
# google - Run account manager as a one-shot task prior to sshd starting.
#
#
start on starting sshd
task

exec /usr/share/google/google_daemon/manage_accounts.py --interval=-1

