# package/upstart/google-ip-forwarding-daemon.conf
# Manages IP forwarding.
start on stopped google-instance-setup

respawn
exec /usr/bin/google_ip_forwarding_daemon
