# package/upstart/google-instance-setup.conf
# Runs instance setup on boot.
start on (started network and started rsyslog) or (started networking)

task

exec /usr/bin/google_instance_setup
