# package/upstart/google-startup-scripts.conf
# google_run_startup_scripts.conf
description "Runs a startup script from metadata."
start on google-instance-setup-ran

exec /usr/bin/google_metadata_script_runner --script-type startup
