# package/upstart/google-ip-forwarding-daemon.conf
# google_ip_forwarding_daemon.conf
description "Manages IP forwarding."
start on google-instance-setup-ran

respawn
exec /usr/bin/google_ip_forwarding_daemon
