# package/upstart/google-shutdown-scripts.conf
# google_run_shutdown_scripts
description "Runs a shutdown script from metadata."
start on starting rc RUNLEVEL=[06]
task

exec /usr/bin/google_metadata_script_runner shutdown
