# package/upstart/google-clock-skew-daemon.conf
# google_clock_skew_daemon.conf
description "Sync the system clock on migration."
start on google-instance-setup-ran

respawn
exec /usr/bin/google_clock_skew_daemon
